# binary/grakn-core.service
[Unit]
Description=Grakn Core
[Service]
Type=simple
ExecStart=/opt/grakn/core/grakn server
